# deployment/coin87-frontend.service
[Unit]
Description=Coin87 Frontend Next.js Service
After=network.target coin87-backend.service
Requires=coin87-backend.service

[Service]
Type=simple
User=coin87
Group=coin87
WorkingDirectory=/home/coin87/coin87sourcev2/frontend
Environment="PATH=/home/coin87/.nvm/versions/node/v18.17.0/bin:/usr/local/bin:/usr/bin:/bin"
Environment="NODE_ENV=production"
Environment="PORT=9011"
EnvironmentFile=/home/coin87/coin87sourcev2/frontend/.env.local

# Build and start
ExecStartPre=/home/coin87/.nvm/versions/node/v18.17.0/bin/npm run build
ExecStart=/home/coin87/.nvm/versions/node/v18.17.0/bin/npm start

# Restart policy
Restart=always
RestartSec=10
StandardOutput=append:/var/log/coin87/frontend.log
StandardError=append:/var/log/coin87/frontend-error.log

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/coin87

[Install]
WantedBy=multi-user.target
