# deployment/coin87-backend.service
[Unit]
Description=Coin87 Backend API Service
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=coin87
Group=coin87
WorkingDirectory=/home/coin87/coin87sourcev2/backend
Environment="PATH=/home/coin87/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/home/coin87/coin87sourcev2/backend"
EnvironmentFile=/home/coin87/coin87sourcev2/backend/.env

# Production command with 4 workers
ExecStart=/home/coin87/.local/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 9010 \
    --workers 4 \
    --log-level info \
    --no-access-log

# Restart policy
Restart=always
RestartSec=10
StandardOutput=append:/var/log/coin87/backend.log
StandardError=append:/var/log/coin87/backend-error.log

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/coin87

[Install]
WantedBy=multi-user.target
