# deployment/coin87-backend.service
[Unit]
Description=Coin87 Backend API Service
After=network.target postgresql.service
Requires=postgresql.service

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/var/www/coin87sourcev2/backend
Environment="PATH=/root/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="PYTHONPATH=/var/www/coin87sourcev2/backend"
EnvironmentFile=/var/www/coin87sourcev2/backend/.env

# Production command with 4 workers
ExecStart=/root/.local/bin/uvicorn app.main:app \
    --host 127.0.0.1 \
    --port 9010 \
    --workers 4 \
    --log-level info \
    --no-access-log

# Restart policy
Restart=always
RestartSec=10
StandardOutput=append:/var/log/coin87/backend.log
StandardError=append:/var/log/coin87/backend-error.log

# Security
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/coin87

[Install]
WantedBy=multi-user.target
